# meshtastic.service
[Unit]
Description=Meshtastic Terminal Monitor with ChatBot
After=network.target

[Service]
Type=simple
User=iain
WorkingDirectory=/home/iain/Meshtastic
ExecStart=/home/iain/Meshtastic/start_terminal.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment settings
Environment="PYTHONUNBUFFERED=1"
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/iain/.Xauthority"

[Install]
WantedBy=multi-user.target
